# deploy/show-watcher.service
[Unit]
Description=Show Watcher dashboard and scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=showwatcher
WorkingDirectory=/opt/movie-tracker
EnvironmentFile=/etc/show-watcher.env
ExecStart=/usr/bin/python3 /opt/movie-tracker/web.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
